# build/lychee-supervisor.ini
# Tell Supervisord not to fork into the background
[supervisord]
nodaemon=true
#user=nginx
pidfile=/control/supervisord.pid

# Override default supervisord socket path, since 
# running as an unprivileged user
#[unix_http_server]
#file=/control/supervisord.sock

# Startup php-fpm7, then nginx
[program:php-fpm7]
command=/usr/sbin/php-fpm7 -c /etc/php7 -F

[program:nginx]
command=/usr/sbin/nginx
